# playbooks/roles/nginx/templates/etc/systemd/system/nginx.service.d/override.conf
[Service]
Restart=always

# default values if we want to change them
# RestartUSec=100ms
# StartLimitBurst=5
# StartLimitInterval=10000000
